# chromeos-base/chromeos-bsp-fydetab_duo-openfyde/files/init/init_leds.conf
# Copyright 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "init system leds"
author        "yang@fydeos.io"

start on started system-services
task
oom score -100

script
  #modprobe -q ledtrig-timer
  #modprobe -q ledtrig-default-on
  udevadm trigger --action=change --subsystem-match=backlight
  echo 55 > /sys/module/cpu_boost/parameters/input_boost_freq_percent
end script
