# chromeos-base/chromeos-bsp-fydetab_duo-openfyde/files/init/init_sdmmc.conf
# Copyright 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "init sdmmc module"
author        "yang@fydeos.io"

start on started system-services
task
oom score -100

exec modprobe -q dw_mmc_rockchip
