# chromeos-base/chromeos-bsp-fydetab_duo-openfyde/files/rules/99-install-himax-mod.rules
ACTION=="add", KERNEL=="mali0", SUBSYSTEM=="misc", RUN+="/sbin/initctl start --quiet --no-wait init_himax_touchscreen"
