# chromeos-base/chromeos-bsp-fydetab_duo-openfyde/files/rules/99-usb-quectel-modem.rules
ACTION=="add", KERNEL=="cdc-wdm0", SUBSYSTEM=="usbmisc", RUN+="/sbin/runuser -u modem -- /usr/bin/quectel-CM"
